# backend-supervisor-csportfolio-sync.conf
[unix_http_server]
file=/var/run/supervisor.sock

[supervisord]
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=www-data

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# PHP-based cron: Run sync-prices.php every hour
# Uses a persistent loop with sleep to avoid supervisor reload issues
[program:csportfolio-sync-prices]
command=/bin/bash -c "while true; do php /var/www/html/api/sync-prices.php; sleep 3600; done"
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/csportfolio-sync.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/csportfolio-sync-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
priority=999
environment=PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
