# supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=***REMOVED***
environment=PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# Start Apache in foreground
[program:apache2]
command=/bin/bash -c "source /etc/apache2/envvars && exec /usr/sbin/apache2ctl -D FOREGROUND"
stdout_logfile=/var/log/supervisor/apache2.log
stderr_logfile=/var/log/supervisor/apache2-error.log
autorestart=true
autostart=true
priority=10
redirect_stderr=true

# PHP-based cron: Run sync-prices.php every hour
[program:cs***REMOVED***-sync-prices]
command=/bin/bash -c "while true; do php /var/www/html/api/sync-prices.php; sleep 3600; done"
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/cs***REMOVED***-sync.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/cs***REMOVED***-sync-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
priority=999
